# xfce4-volumed-pulse/xfce4-volumed-pulse.service
[Unit]
Description=XFCE-Volumed-Pulse
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/xfce4-volumed-pulse
